# stonkbot.service
[Unit]
Description=Stonk Bot

[Service]
ExecStart=/usr/bin/stonkbot
Restart=always

[Install]
WantedBy=default.target

